# distro/airootfs/etc/systemd/system/syslog-ng@default.service.d/no-network-wait.conf
[Unit]
Wants=
After=
After=syslog.socket
