# rootfs/etc/systemd/system/cerberix-update.timer
[Unit]
Description=Cerberix Update Check
[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
[Install]
WantedBy=timers.target
