# rootfs/etc/systemd/system/cerberix-update.service
[Unit]
Description=Cerberix Update Check
[Service]
Type=oneshot
Nice=10
IOSchedulingClass=idle
# This timer is only a best-effort update check. checkupdates exits 1 when a
# mirror or repo fetch fails and 2 when no updates are available; neither
# should leave a systemd failure behind for a background probe.
ExecStart=/bin/sh -c '/usr/bin/checkupdates 2>/dev/null || exit 0'
